# DPU/linux_scripts/rc.local
#! /bin/bash
/etc/init.d/config_after_reboot.sh
exit 0
